# deploy/vm/systemd/bracketiq-verify-host.timer
[Unit]
Description=Run the BracketIQ production VM verification daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=30m
Unit=bracketiq-verify-host.service

[Install]
WantedBy=timers.target
